# docker/supervisor/horizon.conf
[program:horizon]
command=php artisan horizon
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/horizon.log
stderr_logfile=/var/log/supervisor/horizon.err.log
user=www-data
